# server/supervisor/messenger-worker.ini
[supervisord]
nodaemon=true

[include]
files = /etc/supervisor/conf.d/*.conf

[program:cron]
command=/usr/sbin/crond -f -l 8 -L /var/log/cron.log
user=root
autorestart=true

[program:php-fpm]
command=docker-php-entrypoint php-fpm
user = root
autostart = true
autorestart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile = /dev/stderr
stderr_logfile_maxbytes=0

[program:messenger-consume]
command=php /srv/app/bin/console messenger:consume async --time-limit=3600
user=root
numprocs=2
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d
